# s6/service-iplayer-arr/run
#!/command/with-contenv bash
exec s6-setuidgid hotio /app/iplayer-arr
